# timers/lg_renew@.timer
[Unit]
Description=Daily LG dev mode token renew for %i

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
